# trunk/contrib/service/pgsqllogd/run
#!/bin/sh
exec 2>&1
export PATH="$PATH":/usr/local/bin:/usr/local/progs/v-q/bin:/var/vq/bin:/var/v-q/bin
exec pgsqllogd
